# petalinux/meta-daphne/recipes-core/daphne-services/files/daphne.service
[Unit]
Description=DAPHNE Slow Controller
After=hermes.service endpoint.service firmware.service
Requires=endpoint.service firmware.service
PartOf=firmware.service

[Service]
Type=simple
Environment=LD_LIBRARY_PATH=/usr/lib/daphne-server
Environment=DAPHNE_BIND_ENDPOINT=tcp://*:40001
ExecStart=/bin/sh -c 'exec /usr/bin/daphneServer --bind "${DAPHNE_BIND_ENDPOINT}"'
Restart=on-failure
RestartSec=1
TimeoutStopSec=5s

[Install]
WantedBy=multi-user.target
